# nerell-parent/nerell-robot/src/main/scripts/nerell.service
[Unit]
After=multi-user.target

[Service]
ExecStart=/home/pi/nerell/run.sh
WorkingDirectory=/home/pi/nerell

[Install]
WantedBy=default.target